# upstart/authenticatord.conf
description "authenticatord"

start on started m2_node
stop on runlevel [!2345] or rcS or stopped m2_node

respawn

exec /usr/bin/authenticatord
